# usbarmory-control/sys-usb/90-usbarmory.rules
SUBSYSTEM=="net", ACTION=="add", RUN+="/usr/sbin/dhclient $name"
