# usbarmory-control/sys-usb/90-usbarmory.rules
SUBSYSTEM=="net", ACTION=="add", ATTR{address}=="1a:55:89:a2:69:42", RUN+="/usr/sbin/dhclient $name"
