# usbarmory-control/sys-usb/rc.local
#!/bin/sh

ln -s /rw/config/90-usbarmory.rules /etc/udev/rules.d/90-usbarmory.rules
/usr/bin/udevadm control --reload-rules
